# profile.d/01-environment.sh
# Defines default editor
export EDITOR="/bin/vim"
# export EDITOR="mate -w"

# Defines manpager
export MANPAGER='less -Xr'

# Defines my executables
export PATH="$HOME/bin:$PATH"
# export PATH="$HOME:$PATH"