# profile.d/01-nvm.sh
# Defines nvm path
export NVM_DIR="$HOME/.nvm"

# Imports nvm configuration
[ -s "$NVM_DIR/nvm.sh" ] && \. "$NVM_DIR/nvm.sh"  # This loads nvm

# Imports nvm completion
[ -s "$NVM_DIR/bash_completion" ] && \. "$NVM_DIR/bash_completion"  # This loads nvm bash_completion
