# profile.d/02-git-custom.sh
git config --global alias.ci 'commit -S'
git config --global alias.amend "commit -S --amend"
git config --global user.name "Jonathan C. Calixto"
git config --global user.email jonathanccalixto@gmail.com
git config --global user.signingkey "3C15F404B871867F"
