# profile.d/11-oh-my.zsh
# Defines oh-my-zsh plugins to use
# plugins=(git docker rbenv vi-mode rsync)
plugins=(git docker rbenv rsync ruby node)

# Defines oh-my-zsh path as zsh
export ZSH="$HOME/.oh-my-zsh"

# Imports oh-my-zsh configuration
source $ZSH/oh-my-zsh.sh
